# downloadDeleter.timer
[Unit]
Description=Daily stale-download-file deleter

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
